# pundit-agent/system/etc/cron.d/pundit-localization-daemon
*/15 * * * * root perl -I /opt/pundit-agent/lib /opt/pundit-agent/bin/localization_traceroute_daemon.pl &> /opt/pundit-agent/localization.log
